# files/etc/service/dnsmasq/run
#!/bin/execlineb -P
fdmove -c 2 1
dnsmasq -x /dev/null --no-daemon -i wlan0 -z \
  -F ${DHCP_RANGE}
