# files/etc/service/dnsmasq/run
#!/bin/execlineb -P
fdmove -c 2 1
dnsmasq -x '' --no-daemon -i wlan0 -z \
  -F 10.10.0.10,10.10.3.254
