# files/etc/service/kiwix/run
#!/bin/execlineb -P
fdmove -c 2 1
/usr/bin/kiwix-serve --port 8080 --library /data/content/kiwix/library.xml --index /data/content/kiwix/index
