# files/etc/service/kiwix/run
#!/bin/execlineb -P
fdmove -c 2 1
/kiwix/bin/kiwix-serve \
  --library /content/kiwix/library \
  --index   /content/kiwix/index
