# files/etc/service/ka-lite/run
#!/bin/sh
exec 2>&1
HOME=/data/kalite s6-setuidgid kalite /usr/bin/kalite start --foreground
